# spindown.service
[Unit]
Description=Idle disk spindown service

[Service]
ExecStart=/usr/bin/spindownd

[Install]
WantedBy=multi-user.target
